# analysis_report.json
[
  {
    "file": "source_codes\\demo1.py",
    "hash": "de3c00b20c907de986d38b59c63e2c25b840a1bfbdf4bb84781f86fc367150e6",
    "analysis": [
      "**Analysis Summary**",
      "",
      "1. **Purpose & Risk**:  ",
      "   - **Benign** (simulates malicious patterns for analysis/testing). No harmful actions executed.",
      "",
      "2. **Language/Style**:  ",
      "   - Python, cleanly written with comments. Base64 obfuscation for demo purposes (non-malicious).",
      "",
      "3. **Completeness**:  ",
      "   - Self-contained; no external dependencies. Simulated logic (e.g., commented-out registry edits).",
      "",
      "4. **Hardcoded Paths/Names**:  ",
      "   - `./.persistence_demo.txt` (harmless temp file). No suspicious paths.",
      "",
      "5. **File Operations**:  ",
      "   - Creates `.persistence_demo.txt` (non-malicious content).",
      "",
      "6. **Forensic Artifacts**:  ",
      "   - Temp file, printed system info, no network traffic, process logs of benign activity.",
      "",
      "7. **APIs/System Calls**:  ",
      "   - `os`, `platform`, `socket`, `base64`, `pathlib`. No malicious execution.",
      "",
      "8. **Abused APIs**:  ",
      "   - `base64` (common in obfuscation), `socket` (C2 simulation). Not actively abused here.",
      "",
      "9. **Command Execution**:  ",
      "   - None. `ENCODED_CMD` decodes to `ping -c \"test\"` (not executed).",
      "",
      "10. **Suspicious Strings**:  ",
      "    - `http://example.com:8080/api` (benign C2 placeholder), VM check strings.",
      "",
      "11. **Obfuscation**:  ",
      "    - Base64 encoding (simulated, not hiding malicious content).",
      "",
      "12. **Crypto Functions**:  ",
      "    - Base64 decoding (used for demo strings, no encryption).",
      "",
      "13. **Credentials/API Keys**:  ",
      "    - None.",
      "",
      "14. **OS Interaction**:  ",
      "    - Simulated persistence (file write), no registry/services modified.",
      "",
      "15. **Persistence/Priv Esc**:  ",
      "    - Simulated via temp file; no real persistence.",
      "",
      "16. **Security Bypass**:  ",
      "    - Anti-VM check (detects VMs but takes no action).",
      "",
      "17. **Data Collection**:  ",
      "    - Collects system info (cwd, PID, OS) for demo.",
      "",
      "18. **Network Features**:  ",
      "    - Prepares socket for `example.com:8080` (no connection made).",
      "",
      "19. **Protocols**:  ",
      "    - HTTP (simulated, no traffic).",
      "",
      "20. **C2/Exfiltration**:  ",
      "    - C2 URL decoded but unused. No data exfiltration.",
      "",
      "21. **Lateral Movement**:  ",
      "    - None.",
      "",
      "22. **File Operations**:  ",
      "    - Writes harmless text file.",
      "",
      "23. **Data Theft**:  ",
      "    - None.",
      "",
      "24. **Exfiltration Prep**:  ",
      "    - None.",
      "",
      "25. **Anti-Analysis**:  ",
      "    - Basic VM check (no evasion beyond detection).",
      "",
      "26. **Delays/Sandbox Evasion**:  ",
      "    - None.",
      "",
      "27. **Polymorphism**:  ",
      "    - None.",
      "",
      "28. **Suspicious Lines**:  ",
      "    - **Lines 10-12**: Base64-encoded C2/commands (benign but mimics malware).  ",
      "    - **Lines 16-20**: Anti-VM check (common in malware).  ",
      "    - **Lines 30-40**: Simulated persistence logic (non-functional).",
      "",
      "29. **Dead Code**:  ",
      "    - Commented-out registry edits (simulated intent).",
      "",
      "30. **MITRE ATT&CK**:  ",
      "    - **T1083** (File Discovery), **T1497** (Virtualization/Sandbox Evasion), **T1547** (Persistence), **T1583** (C2 Infrastructure).",
      "",
      "31. **Impact**:  ",
      "    - None (code is inert). Simulated actions lack real-world effect.",
      "",
      "32. **Risk Rating**:  ",
      "    - **Low** (no harmful execution).",
      "",
      "33. **Next Steps**:  ",
      "    - Confirm inertness via sandbox. Check for live C2`example.com` (unlikely). Monitor for similar patterns in real code.",
      "",
      "34. **Unclear Areas**:  ",
      "    - Simulated functions (e.g., `connect_c2_demo`) require dynamic analysis to confirm no hidden behavior.",
      "",
      "35. **Summary**:  ",
      "    - **Benign demo code** mimicking malware TTPs (C2 obfuscation, VM checks, persistence).  ",
      "    - No active payloads, network calls, or system modifications.  ",
      "    - Low risk; useful for testing static analysis tools."
    ]
  },
  {
    "file": "source_codes\\demo2.ps1",
    "hash": "8538cc55a45a97ccc727bc296417602733cf1767efe934170f4bf7fc646348c3",
    "analysis": [
      "**Analysis Summary:**",
      "",
      "1. **Purpose & Risk:**  ",
      "   - **Benign** demonstration script for static analysis training. Contains inactive malware-like patterns (commented-out persistence, anti-VM checks) but no harmful execution.",
      "",
      "2. **Code Style:**  ",
      "   - **Obfuscation Demo:** Base64-encoded URL (`http://example.com:8080/heartbeat`), but decoded for clarity. No advanced obfuscation.",
      "",
      "3. **Completeness:**  ",
      "   - Self-contained; no external dependencies.",
      "",
      "4. **Hardcoded Paths/Artifacts:**  ",
      "   - Creates `.persistence_demo.txt` in the script directory (harmless).  ",
      "   - **Registry Key Example** (commented): `HKCU\\...\\Run\\DemoApp`.",
      "",
      "5. **Forensic Artifacts:**  ",
      "   - `.persistence_demo.txt` file.  ",
      "   - Process names logged (VM detection).  ",
      "   - **No network traffic** (URL parsed but not contacted).",
      "",
      "6. **APIs/System Calls:**  ",
      "   - `Get-Process` (VM detection), `Out-File`, Base64 encoding.",
      "",
      "7. **Anti-Analysis:**  ",
      "   - **Anti-VM checks** for `vboxservice`/`vmtoolsd` (T1497).  ",
      "   - **Commented persistence** (T1547.001).",
      "",
      "8. **Suspicious Strings:**  ",
      "   - `example.com:8080` (benign, no C2 activity).  ",
      "   - Base64-encoded URL (decoded at runtime).",
      "",
      "9. **MITRE ATT&CK:**  ",
      "   - **T1497** (Virtualization Evasion), **T1053** (Scheduled Task – if uncommented), **T1071.001** (HTTP for C2 – demo only).",
      "",
      "10. **Impact/Risk:**  ",
      "    - **Low risk** – No active malicious payload. Demonstrates techniques without execution.",
      "",
      "11. **Next Steps:**  ",
      "    - **Static Analysis:** Confirm no uncommented dangerous code.  ",
      "    - **Sandbox Testing:** Verify no network calls/file changes.  ",
      "    - **Monitor:** Registry writes if script is modified.",
      "",
      "**Key Suspicious Lines:**  ",
      "- **Lines 4-5:** Base64 URL encoding/decoding (common in malware).  ",
      "- **Lines 8-14:** Anti-VM checks (T1497).  ",
      "- **Line 20:** Commented registry persistence (T1547.001).  ",
      "",
      "**Conclusion:**  ",
      "Training script mimicking malware techniques without harmful intent. Risk remains **Low** unless modified."
    ]
  }
]